# systemd/xray.service.d/99-tun-routing.conf
[Service]
# 必须显式声明，否则脚本可能找不到 ip 命令或者无法操作网络
# AmbientCapabilities 已经在主服务里声明了，这里会自动继承

# 告诉脚本配置文件在哪里（虽然主服务有，但为了保险再次明确）
Environment=CONFIGURATION_DIRECTORY=/etc/xray

# 启动后执行脚本
# 这里的 - 是为了让 systemd 忽略错误（万一脚本挂了不影响 Xray 运行，方便调试）
# 但生产环境建议去掉 -，确保路由配置成功
ExecStartPost=/usr/local/bin/xray-tun-helper.sh start

# 停止后执行清理
# 这一步非常重要！否则你停止服务后，路由规则还在，会导致断网（因为 xray0 没了）
ExecStopPost=/usr/local/bin/xray-tun-helper.sh stop
